# server-conf/supervisor-nginx.conf
[program:nginx]
command = service nginx start
stdout_events_enabled=true
stderr_events_enabled=true
