# server-conf/supervisor-uwsgi.conf
[program:uwsgi]
command = /usr/bin/uwsgi --ini /etc/uwsgi/apps-enabled/default

